# packaging/deb/prerm
#!/bin/sh
set -e

systemctl stop rfmpd 2>/dev/null || true
systemctl disable rfmpd 2>/dev/null || true
